# deploy/facturas-web.service
[Unit]
Description=Facturas Web (FastAPI)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/facturas-app/generar_para_email
Environment=PYTHONUNBUFFERED=1
Environment=WEB_SESSION_SECRET=CAMBIA_ESTA_CLAVE_SUPER_SECRETA
Environment=WEB_SESSION_HTTPS_ONLY=CAMBIA_HTTPS_ONLY
ExecStart=/usr/local/bin/uv run uvicorn web.app:app --host 127.0.0.1 --port 8000 --proxy-headers
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
